# cec.cron
# GitHub Update Once A Day
@reboot pi bash /home/pi/scripts/cec/ghupdate.sh &
0 1 * * * pi bash /home/pi/scripts/cec/ghupdate.sh &

# Download calendar and import weekly (based on day selected in system)
5 4 * * * pi curl http://localhost/cec/importentries.php

# Cleanup Once a Week
30 1 * * 2 pi bash /home/pi/scripts/cleanup.sh &
